# ansible/roles/etcd/files/etcd.service
[Unit]
Description=etcd key-value store
Documentation=https://github.com/etcd-io/etcd
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=notify
User=etcd
Group=etcd
ExecStart=/usr/local/bin/etcd
EnvironmentFile=/etc/etcd/etcd.conf
Restart=on-failure
LimitNOFILE=40000
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=120
KillMode=process
SendSIGKILL=no
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/etcd
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
StandardOutput=journal
StandardError=journal
SyslogIdentifier=etcd

[Install]
WantedBy=multi-user.target
